# Services/roaring-mic-routesd.service
[Unit]
Description=Roaring mic router (SM7B/Astro -> B1/B2)
After=pipewire.service pipewire-pulse.service
Wants=pipewire.service pipewire-pulse.service

[Service]
Type=simple
ExecStart=%h/bin/roaring_mic_routesd.sh
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
