# Services/roaring-audio-stackd.service
[Unit]
Description=Roaring audio stack daemon (VM sinks + loopbacks to Astro)
After=pipewire-pulse.service wireplumber.service
Wants=pipewire-pulse.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/bin/roaring_audio_stackd.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
